# .docker-compose/web/etc/cron.d/root
* * * * *   (cd /var/www/html && echo "running ...")
